# .config/systemd/user/git-auto-sync.service
[Unit]
Description=Automated Git Sync Service

[Service]
Type=oneshot
ExecStart=%h/scripts/git-auto-sync.sh
