# .config/autostart/vesktop.desktop
[Desktop Entry]
Type=Application
Name=Vesktop
Comment=Vesktop autostart script
Exec=/opt/Vesktop/vesktop --relaunch --start-minimized
StartupNotify=false
Terminal=false
Icon=vesktop